# image_builder/data/jammy-rt-humble-stanley/rootfs/etc/profile.d/99-source-ros.sh
if [ -f /opt/ros/humble/setup.bash ]; then
  source /opt/ros/humble/setup.bash
fi
if [ -f /opt/pixx_ws/install_aarch64/setup.bash ]; then
  source /opt/pixx_ws/install_aarch64/local_setup.bash
fi