# image_builder/data/jammy-rt-humble-stanley/rootfs/etc/profile.d/99-source-ros.sh
if [ -d "/opt/ros/galactic" ]; then
  source /opt/ros/galactic/setup.bash
fi
if [ -d "/opt/ros/humble" ]; then
  source /opt/ros/humble/setup.bash
fi
