# image_builder/data/jammy-rt-humble-stanley/rootfs/etc/cron.d/time_ntp_job
*/10 * * * * pix /opt/pixx_ws/pixx_scripts/time-ntp
