# image_builder/data/jammy-rt-humble-stanley/rootfs/etc/cron.d/time_ntp_job
*/10 * * * * stan /opt/stanley2/scripts/time-ntp
